# systemd/voice-codex.service
[Unit]
Description=Voice Codex desktop assistant
After=graphical-session.target

[Service]
Type=simple
Environment=VOICE_CODEX_START_HIDDEN=1
ExecStart=/home/juanma/AIComputerbyAmccgil/scripts/voice-codex
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
